# roles/backup-lima-city-trigger/templates/lima-city-backup.timer.j2
[Unit]
Description=Daily lima-city backup trigger timer

[Timer]
OnCalendar={{ lima_city_trigger_schedule }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
