# roles/backup-sync/templates/backup-check.timer.j2
[Unit]
Description=Weekly backup integrity check timer

[Timer]
OnCalendar={{ backup_check_schedule }}
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
